# .config/systemd/user/hypr-name-auto.service
[Unit]
Description=Auto-rename hyprland workspaces
After=graphical-session.target

[Service]
Type=simple
Environment="XDG_RUNTIME_DIR=/run/user/%U"
Environment="DISPLAY=:0"
ExecStartPre=/usr/bin/bash -c 'for i in {1..30}; do \
    [ -S "$XDG_RUNTIME_DIR/hyprland/hyprland.sock" ] && exit 0; \
    sleep 1; \
  done; \
  echo "Timed out waiting for Hyprland socket" >&2; \
  exit 1'
ExecStart=%h/.local/bin/hypr-configs/hypr-name-auto.sh
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
