# .config/systemd/user/rclone-bisync.service
[Unit]
Description=Sync Google Drive

[Service]
Type=oneshot
ExecStart=rclone bisync gdrive: ~/Drive-Google
