# .config/systemd/user/rclone-bisync.service
[Unit]
Description=Sync Google Drive

[Service]
Type=oneshot
ExecStart=/usr/bin/rclone bisync gdrive: Drive-Google
